# etc/boatd.service
[Unit]
Description=Robotic Sailing Boat Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/boatd

[Install]
WantedBy=default.target
